# deploy/hantu.service
[Unit]
Description=Hantu Quant Trading System
Documentation=https://github.com/This-HW/hantu_quant
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/hantu_quant

# 시작 시 docker-compose up
ExecStart=/usr/bin/docker compose up -d
# 중지 시 docker-compose down
ExecStop=/usr/bin/docker compose down
# 재시작 시 recreate
ExecReload=/usr/bin/docker compose up -d --force-recreate

# 타임아웃 설정
TimeoutStartSec=120
TimeoutStopSec=30

# 재시작 정책
Restart=on-failure
RestartSec=10

# 로깅
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hantu

[Install]
WantedBy=multi-user.target
